# deploy/systemd/drl-trading-ui.service
[Unit]
Description=DRL Trading System - Streamlit UI
After=network.target drl-trading-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/projects/drl-trading-system/repo
EnvironmentFile=/root/.openclaw/projects/drl-trading-system/repo/.env
Environment=API_SERVER_URL=http://127.0.0.1:5001
ExecStart=/usr/bin/python3 -m streamlit run src/ui/app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --server.enableCORS false --server.enableXsrfProtection false --browser.gatherUsageStats false
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
